# crontab
*/1 * * * * python /task.py >> /logfile
